# openwrt/files/etc/crontabs/root
0 */1 * * * /usr/sbin/ntpdate ntp.tencent.com && /usr/sbin/hwclock -w
