# openwrt/files/etc/crontabs/root
0 */1 * * * /usr/sbin/ntpdate time1.cloud.tencent.com && /usr/sbin/hwclock -w
